# supervisor/conf.d/kibana.conf
[program:kibana]
command=/opt/kibana/bin/kibana
stdout_logfile	= /var/log/supervisor/%(program_name)s.log
stderr_logfile	= /var/log/supervisor/%(program_name)s.log
